# besd/supervisord.conf
[program:besd]
command=/usr/bin/besdaemon -i /usr -c /etc/bes/bes.conf -r /var/run/bes
